# profile.d/chef.sh
export CHEF_SERVER=prod
export CHEF_USER=geoff
