# profile.d/history.sh
export HISTTIMEFORMAT="%F %T "
export HISTSIZE=5000
