# profile.d/vertx.sh
if [ -d "/opt/vert.x" ] ; then
  export VERTX_HOME=/opt/vert.x
  export PATH=$VERTX_HOME/bin:$PATH
fi
